# systemd/erebine-eem-agent.service
# SPDX-License-Identifier: MIT
# Erebine - XEM Agent systemd unit (installed by the erebine-eem-agent
# package). Configure /etc/erebine/eem-agent.env, then:
#   sudo systemctl enable --now erebine-eem-agent

[Unit]
Description=Erebine XEM Agent
Documentation=https://erebine.ai/docs/private-agents
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=erebine
Group=erebine

ExecStart=/usr/bin/erebine-eem-agent

# Environment configuration
Environment=LOG_LEVEL=info
Environment=TOOLS_DIR=/var/lib/erebine/tools
EnvironmentFile=-/etc/erebine/eem-agent.env

# Runtime directories
RuntimeDirectory=erebine-eem-agent
StateDirectory=erebine

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/erebine
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
